# supervisord.conf
[supervisord]
nodaemon=true

[program:gwd]
command=/geneweb/distribution/gwd

[program:gwsetup]
command=/geneweb/distribution/gwsetup -lang fr